# .config/systemd/user/wallpaper-daily.timer
[Unit]
Description=Daily wallpaper rotation timer

[Timer]
OnCalendar=daily
Persistent=true
Unit=wallpaper-daily.service

[Install]
WantedBy=timers.target
